# MSG-WR211_CMCC/source/vendors/Broadcom/bcm947xx/inittab
ttyS0::respawn:/bin/sh
::sysinit:/etc_ro/rcS
